# systemd/sysstat-report@.timer
[Unit]
Description=Send %i Sysstat email reports

[Timer]
OnCalendar=%i
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
